# dockerlab/templates/systemd/dockerlab-deploy.timer
[Unit]
Description=Run Docker Git Deploy periodically

[Timer]
OnBootSec=1min
OnUnitInactiveSec=5min

[Install]
WantedBy=timers.target
