# docker-git-deploy-skill/assets/systemd/docker-git-deploy.timer.in
[Unit]
Description=Run Docker Git Deploy periodically

[Timer]
OnBootSec=1min
# @@POLL_INTERVAL@@ is substituted by install.sh (e.g. 5min).
OnUnitInactiveSec=@@POLL_INTERVAL@@
Persistent=true

[Install]
WantedBy=timers.target
